# picoclaw/rootfs/etc/services.d/terminal/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

exec ttyd -W -p 17681 -t titleFixed="PicoClaw Terminal" /usr/bin/picoclaw-shell
