# picoclaw/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

mkdir -p /run/nginx
exec nginx -g 'daemon off;'
