# picoclaw/rootfs/etc/services.d/webui/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

exec /usr/bin/picoclaw-control-panel
